# SRC-GNB-PLFM-B0-v1.5.0-RC2.1/buildroot/board/edgecompute/raptor2/rootfs-overlay/etc/systemd/system/dpdk.service
[Unit]
Description=Service to setup DPDK hugepages environment

[Service]
User=root
ExecStart=/usr/bin/python3 /usr/bin/dpdk-hugepages.py -p 32M --setup 512M
Restart=no

[Install]
WantedBy=multi-user.target
