# SRC-GNB-PLFM-B0-v1.5.0-RC2.1/buildroot/board/edgecompute/raptor2/rootfs-overlay/etc/systemd/system/logger.service
[Unit]
Description=Service to Launch Raptor2 Logger Service

[Service]
User=root
WorkingDirectory=/root
ExecStart=/bin/bash '/etc/l1logger_service.sh'
Restart=always

[Install]
WantedBy=multi-user.target
