# SRC-GNB-PLFM-B0-v1.5.0-RC2.1/buildroot/board/edgecompute/raptor2/rootfs-overlay/etc/systemd/system/monit.service
[Unit]
Description=Service to Launch Monit Service

[Service]
Type=forking
User=root
ExecStart=/bin/monit
RemainAfterExit=true
Restart=on-failure

[Install]
WantedBy=multi-user.target
