# SRC-GNB-PLFM-B0-v1.5.0-RC2.1/buildroot/board/edgecompute/raptor2/rootfs-overlay/etc/systemd/system/platform.service
[Unit]
Description=Service to setup Platform Specific environment

[Service]
User=root
ExecStart=/bin/bash '/etc/setup_platform.sh'
Restart=no

[Install]
WantedBy=multi-user.target
