# SRC-GNB-PLFM-B0-v1.5.0-RC2.1/buildroot/board/edgecompute/raptor2/rootfs-overlay/etc/systemd/system/system_logger.service
[Unit]
Description=Service to Launch Raptor2 System logger Service
Before=monit.service

[Service]
User=root
WorkingDirectory=/root
Type=oneshot
ExecStart=/bin/bash '/etc/system_logger.sh'
Restart=on-failure
 
[Install]
WantedBy=multi-user.target
