# SRC-GNB-PLFM-B0-v1.5.0-RC2.1/buildroot/board/edgecompute/raptor2/rootfs-overlay/etc/systemd/system/timesync.service
[Unit]
Description=Service to Launch Raptor2 Time Sync Service
Before=raptor2.service
After=network-online.target
Requires=network-online.target

[Service]
User=root
WorkingDirectory=/root
ExecStart=python3 '/etc/timesync_service.py'
Restart=always

[Install]
WantedBy=multi-user.target
